# tasks/celery_tasks.conf
[program:celery_tasks]
command=/usr/local/bin/celery worker --workdir=/opt/python/rui/ -A tasks -l info -Q ck:celery:mail
autostart=true
autorestart=true
stdout_logfile=/var/log/celery_tasks.log